# roles/maintenance/templates/mpserver-maintenance.service.j2
[Unit]
Description=server maintenance

[Service]
Type=oneshot
ExecStart=/home/{{ username }}/maintenance/mpserver-maintenance.sh
